# crontab
* * * * * root /bin/sh /scripts/hello-world.sh
